# systemd/home-chen-smb.automount
[Unit]
Description=Automount vps-smb

[Automount]
Where=/home/chen/smb

[Install]
WantedBy=multi-user.target
